# task-2.1/crontab
0 0 * * * find /tmp/* -type f -ctime +14 -size +5M -exec rm -f {} \;
